# example/bash/.deploy/supervisord.conf
[program:APP_NAME]
command=DEPLOY_DIR/SERVICE_COMMAND
process_name=%(program_name)s
numprocs=1
autostart=true
autorestart=true
user=deploygirlexample
environment=SERVICE_ENV
